# rootfs/etc/services.d/growmind/run
#!/usr/bin/with-contenv bashio
set -euo pipefail

export GEMINI_API_KEY="$(bashio::config 'gemini_api_key')"
export GEMINI_MODEL="$(bashio::config 'gemini_model')"
export LOG_LEVEL="$(bashio::config 'log_level')"
export HASS_TOKEN="$(bashio::config 'hass_token')"
export INFLUX_URL="$(bashio::config 'influx_url')"
export INFLUX_TOKEN="$(bashio::config 'influx_token')"
export INFLUX_ORG="$(bashio::config 'influx_org')"
export INFLUX_BUCKET="$(bashio::config 'influx_bucket')"
export GRAFANA_EMBED_URL="$(bashio::config 'grafana_embed_url')"

if [ -z "${GEMINI_API_KEY:-}" ]; then
  bashio::log.warning "GEMINI_API_KEY not set – AI features will be unavailable."
fi

if [ -n "${HASS_TOKEN:-}" ]; then
  bashio::log.info "HASS_TOKEN configured via add-on options."
fi

cd /app/backend
exec uvicorn app.main:app --host 0.0.0.0 --port 8080
